# etc/profile.d/hadoop.sh
export HADOOP_HOME=/home/student/hadoop/hadoop-2.7.2
export PATH=$PATH:${HADOOP_HOME}/bin:${HADOOP_HOME}/sbin
export HADOOP_CONF_DIR=${HADOOP_HOME}/etc/hadoop