# etc/profile.d/hosts.sh
cat > /etc/hosts << EOF
127.0.0.1 localhost
222.27.165.206 hadoop1
222.27.165.205 hadoop3
222.27.165.207 hadoop2
EOF
